# install/proxops.service
[Unit]
Description=ProxOps — GitOps VM lifecycle manager for Proxmox
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/proxops -config /opt/proxops/config.yml
Restart=always
RestartSec=10
WorkingDirectory=/opt/proxops

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=proxops

[Install]
WantedBy=multi-user.target
